# post_install/files/systemd/system/update-issue.timer
[Unit]
Description=Timer for update-issue.service

[Timer]
OnCalendar=*-*-* *:*:00
OnBootSec=0sec
Persistent=true

[Install]
WantedBy=timers.target
