# post_install/files/systemd/system/update-issue.service
[Unit]
Description=Automatically update /etc/issue

[Service]
Type=oneshot
ExecStart=/etc/update_issue.d/update_issue.sh
